# gaea/s6-overlay/s6-rc.d/nginx-svc/run
#!/bin/sh
# nginx-svc (longrun): 前端静态资源 + API 反代
exec nginx -g "daemon off;"
